# php/php74/conf/supervisord.conf
[supervisord]
#后台运行
nodaemon=false
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php]
# variables_order=EGPCS 指定php$_ENV依次从Environment,GET,Post,Cookies取值
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan serve --host=0.0.0.0 --port=80
user=www-data
environment=LARAVEL_WWW_DATA="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
